# docker/s6-overlay/s6-rc.d/atd/run
#!/command/execlineb -P
# at daemon — runs scheduled firmware upgrades
/usr/sbin/atd -f
