# docker/s6-overlay/s6-rc.d/cmds2-api/run
#!/command/execlineb -P

# CMDS2 FastAPI backend (uvicorn)
# Listens on 127.0.0.1:8000, proxied by nginx

cd /root
exec uvicorn web.backend.main:app --host 127.0.0.1 --port 8000 --workers 1 --log-level info
