# docker/s6-overlay/s6-rc.d/init-ssl/up
#!/command/execlineb -P

# Generate self-signed SSL certificate if not present
foreground { /bin/bash /docker-generate-ssl-cert.sh }
