# docker/s6-overlay/s6-rc.d/tftpd/run
#!/command/execlineb -P
# TFTP server — serves firmware images and receives config backups
/usr/sbin/in.tftpd -L -c -p -s /var/lib/tftpboot --foreground
