# docker/s6-overlay/s6-rc.d/httpd/run
#!/command/execlineb -P
# Apache httpd — serves firmware images via HTTP for switch downloads
/usr/sbin/httpd -DFOREGROUND
